# etc/supervisord/devcenter.conf
[program:node]
command = node app.js
directory = /src/DevCenter.me
